# manjaro-i3/etc/profile.d/intel-media.sh
# add libva driver path/name exporting for intel media solution
# uncomment the LIBVA_DRIVER_NAME line to use the Intel Media Driver (iHD) for VAAPI
export LIBVA_DRIVERS_PATH=/usr/lib/dri
export LIBVA_DRIVER_NAME=iHD
